# luci-app-dhcp-comment/files/etc/uci-defaults/99-luci-app-dhcp-comment
#!/bin/sh

# 对原生 Network DHCP 页面与 Status 概览页面实施差量增强注入，百分百保留原版所有功能与布局
python3 -c '
import os

dhcp_path = "/www/luci-static/resources/view/network/dhcp.js"
if os.path.exists(dhcp_path):
    with open(dhcp_path, "r", encoding="utf-8") as f: code = f.read()

    # 1. 静态分配插入单行 Comment 控件
    if "var co = ss.option" not in code:
        t1 = "so = ss.option(form.Value, \x27leasetime\x27,"
        r1 = "var co = ss.option(form.Value, \x27comment\x27, _(\x27Comment\x27));\n\t\tco.rmempty = true;\n\n\t\t" + t1
        code = code.replace(t1, r1)

    # 2. 关联显示活动租约中文备注
    if "mac_comments" not in code:
        t2 = "return callDHCPLeases().then(function(leaseinfo) {"
        r2 = "return Promise.all([ callDHCPLeases(), L.resolveDefault(uci.load(\x27dhcp\x27)) ]).then(function(res) {\n\t\t\t\t\tconst leaseinfo = res[0] || {};\n\t\t\t\t\tconst mac_comments = {}; uci.sections(\x27dhcp\x27, \x27host\x27).forEach(function(s) { L.toArray(s.mac).forEach(function(m) { if (s.comment) mac_comments[m.toLowerCase()] = s.comment; }); });"
        code = code.replace(t2, r2)

        code = code.replace(
            "const columns = [\n\t\t\t\t\t\t\t\t\x27%h\x27.format(host || \x27-\x27),",
            "var cmt = lease.macaddr ? mac_comments[lease.macaddr.toLowerCase()] : null;\n\t\t\t\t\t\t\tif (cmt) { var raw_h = lease.hostname || name; host = raw_h ? (cmt + \x27 (\x27 + raw_h + \x27)\x27) : cmt; }\n\n\t\t\t\t\t\t\tconst columns = [\n\t\t\t\t\t\t\t\t\x27%s\x27.format(host || \x27-\x27),"
        )

        code = code.replace(
            "const columns = [\n\t\t\t\t\t\t\t\t\x27%h\x27.format(host || \x27-\x27),\n\t\t\t\t\t\t\t\tlease.ip6addrs",
            "var cmt6 = lease.macaddr ? mac_comments[lease.macaddr.toLowerCase()] : null;\n\t\t\t\t\t\t\tif (cmt6) { var raw_h = lease.hostname || name; host = raw_h ? (cmt6 + \x27 (\x27 + raw_h + \x27)\x27) : cmt6; }\n\n\t\t\t\t\t\t\tconst columns = [\n\t\t\t\t\t\t\t\t\x27%s\x27.format(host || \x27-\x27),\n\t\t\t\t\t\t\t\tlease.ip6addrs"
        )

    # 3. 边界防御防护
    if "Array.isArray(host.macaddrs)" not in code:
        code = code.replace(
            "o.value(host.macaddrs[0],",
            "if (host && Array.isArray(host.macaddrs) && host.macaddrs.length > 0) o.value(host.macaddrs[0],"
        )

    with open(dhcp_path, "w", encoding="utf-8") as f: f.write(code)

status_path = "/www/luci-static/resources/view/status/include/40_dhcp.js"
if os.path.exists(status_path):
    with open(status_path, "r", encoding="utf-8") as f: code = f.read()

    if "uci.load(\x27dhcp\x27)" not in code:
        code = code.replace("network.getHostHints()", "network.getHostHints(),\n\t\t\tL.resolveDefault(uci.load(\x27dhcp\x27))")

    if "mac_comments" not in code:
        t = "const machints = host_hints.getMACHints(false);"
        h = "const machints = host_hints.getMACHints(false);\n\t\tvar mac_comments = {}; uci.sections(\x27dhcp\x27, \x27host\x27).forEach(function(s) { L.toArray(s.mac).forEach(function(m) { if (s.comment) mac_comments[m.toLowerCase()] = s.comment; }); });\n\t\t"
        code = code.replace(t, h)

        code = code.replace(
            "if (hint && lease.hostname && lease.hostname != hint[1])",
            "if (lease.macaddr && mac_comments[lease.macaddr.toLowerCase()]) { var cmt = mac_comments[lease.macaddr.toLowerCase()], raw_h = lease.hostname || (hint ? hint[1] : null); host = raw_h ? (cmt + \x27 (\x27 + raw_h + \x27)\x27) : cmt; }\n\t\t\telse if (hint && lease.hostname && lease.hostname != hint[1])"
        )
        code = code.replace("\x27%h\x27.format(host || \x27-\x27)", "\x27%s\x27.format(host || \x27-\x27)")

    with open(status_path, "w", encoding="utf-8") as f: f.write(code)
'

rm -rf /tmp/luci-modulecache/* /tmp/luci-indexcache 2>/dev/null || true
exit 0
